# docker/root-filesystem/etc/profile.d/clearpath-ros-environment.sh
# shellcheck shell=bash disable=SC1091
source /etc/ros/setup.bash
